# scripts/systemd/claude-relay-logrotate.timer
[Unit]
Description=Check Claude Relay Service logs for rotation hourly

[Timer]
OnCalendar=hourly
Persistent=true
AccuracySec=1min
Unit=claude-relay-logrotate.service

[Install]
WantedBy=timers.target
